# Agents/Linux/eagle_agent.service
[Unit]
Description=Eagle Agent Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/eagle_agent
Restart=always
User=root

[Install]
WantedBy=multi-user.target